# docker/scheduler/crontab
# PMPulse Laravel Scheduler
# Run the Laravel scheduler every minute
* * * * * cd /var/www/html && php artisan schedule:run >> /var/www/html/storage/logs/scheduler.log 2>&1
